# client/camera.service
[Unit]
Description=Camera client


[Service]
ExecStart=/home/pg/start.sh 100.86.75.91 12345
Restart=on-failure

[Install]
WantedBy=default.target
