# client/camera.service
[Unit]
Description=Camera client


[Service]
ExecStart=/home/pg/start.sh 192.168.1.129 12345
Restart=on-failure

[Install]
WantedBy=default.target
